# msc-wis2node-management/docker/msc-wis2node-management.cron
0 0 * * * msc-wis2node metrics collect-metrics
