# msc-wis2node-management/docker/msc-wis2node-management.cron
# get metrics daily at 23:59Z
59 23 * * * msc-wis2node metrics get
# delete metrics daily at 0Z
0 0 * * * msc-wis2node metrics delete
